# Task_9/files/showlog.service
[Unit]
Description=My showlog service

[Service]
Type=oneshot
EnvironmentFile=/etc/default/showlog.conf
ExecStart=/opt/showlog.sh $WORD $LOG